# com.competitor.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.competitor.monitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/rena/competitor_monitor/venv/bin/python</string>
        <string>/Users/rena/competitor_monitor/run_scrape.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/rena/competitor_monitor</string>
    <key>StartCalendarInterval</key>
    <array>
        <!-- 1st of month at 9am -->
        <dict>
            <key>Day</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- 15th of month at 9am -->
        <dict>
            <key>Day</key>
            <integer>15</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/rena/competitor_monitor/scrape.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rena/competitor_monitor/scrape_error.log</string>
</dict>
</plist>
